# conf/supervisor.ini
[fcgi-program:{{ project_name }}]
socket=tcp://localhost:10234
command=/opt/{{ project_name }}/venv/bin/daphne -u /opt/data/{{ project_name }}/{{ project_name }}%(process_num)d.sock --fd 0 --access-log - --proxy-headers {{ project_name }}.asgi:application
directory=/opt/{{ project_name }}
numprocs=4
process_name={{ project_name }}%(process_num)d
user=root
autostart=true
autorestart=false
stdout_logfile=/opt/{{ project_name }}/logs/{{ project_name }}.log
redirect_stderr=true

[program:celery_beat]
command=/opt/{{ project_name }}/venv/bin/celery -A {{ project_name }} beat -l info
directory=/opt/{{ project_name }}
user=root
autostart=true
autorestart=true
stdout_logfile=/opt/{{ project_name }}/logs/celery_beat.log

[program:celery_worker]
command=/opt/{{ project_name }}/venv/bin/celery -A {{ project_name }} worker -l info
directory=/opt/{{ project_name }}
user=root
autostart=true
autorestart=true
stdout_logfile=/opt/{{ project_name }}/logs/celery_worker.log

